# laravelWorkers/laravel-worker1.conf
[program:laravel-worker1]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan start-added-video-consumer
autostart=true
autorestart=true
user=root
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/supervisor_queue-work1.log


